# docker_app/debian_8_jessie/php/7.0/supervisord.conf
[supervisord]
nodaemon=true

[program:phpfpm]
command=/bin/bash -c "/usr/sbin/php-fpm --fpm-config /etc/php/php-fpm.conf -F"
